# hvac-pygame/hvac-pygame.service
[Unit]
Description=HVAC-Vibe pygame Display
After=network.target mosquitto.service hvac-engine.service
Wants=mosquitto.service hvac-engine.service

[Service]
Type=simple
User=mitkatch
WorkingDirectory=/home/mitkatch/hvac-pygame
Environment=SDL_VIDEODRIVER=fbcon
Environment=SDL_FBDEV=/dev/fb1
ExecStart=/home/mitkatch/hvac-pygame/venv/bin/python main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hvac-pygame

[Install]
WantedBy=multi-user.target
