# hvac-display/hvac-display.service
[Unit]
Description=HVAC-Vibe Display Server
After=network.target mosquitto.service hvac-engine.service
Wants=mosquitto.service hvac-engine.service

[Service]
Type=simple
User=mitkatch
WorkingDirectory=/home/mitkatch/hvac-display
ExecStart=/home/mitkatch/hvac-display/venv/bin/python display_main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hvac-display

[Install]
WantedBy=multi-user.target
